# build-pkg/clixon-vpp_1.1.0_amd64/etc/systemd/system/clixon-vpp-restconf.service
[Unit]
Description=Clixon VPP RESTCONF Daemon
Documentation=https://clixon-docs.readthedocs.io/
After=network.target clixon-vpp-backend.service
Requires=clixon-vpp-backend.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/clixon_restconf -f /etc/clixon/clixon-vpp.xml
Restart=on-failure
RestartSec=5
User=www-data
Group=clixon

[Install]
WantedBy=multi-user.target
