# build-pkg/clixon-vpp_1.1.0_amd64/etc/systemd/system/vpp.service.d/restart-clixon.conf
[Service]
# Restart Clixon in background after VPP starts (wait 90s for full VPP init)
ExecStartPost=-/bin/bash -c '(sleep 90 && systemctl restart clixon-vpp-backend) &'
