# debian/clixon-vpp-backend.service
[Unit]
Description=Clixon VPP Backend Daemon
Documentation=https://clixon-docs.readthedocs.io/
After=network.target vpp.service
Wants=vpp.service

[Service]
Type=forking
Environment="LD_LIBRARY_PATH=/opt/clixon-vpp/lib"
ExecStart=/opt/clixon-vpp/sbin/clixon_backend -f /etc/clixon/clixon-vpp.xml -s running
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/var/run/clixon-vpp.pid
Restart=on-failure
RestartSec=5
User=root
Group=clixon
PermissionsStartOnly=true

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
