# debian/clixon-vpp-restconf.service
[Unit]
Description=Clixon VPP RESTCONF Daemon
Documentation=https://clixon-docs.readthedocs.io/
After=network.target clixon-vpp-backend.service
Requires=clixon-vpp-backend.service

[Service]
Type=simple
Environment="LD_LIBRARY_PATH=/opt/clixon-vpp/lib"
ExecStart=/opt/clixon-vpp/sbin/clixon_restconf -f /etc/clixon/clixon-vpp.xml -D 1
Restart=on-failure
RestartSec=5
User=www-data
Group=clixon

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/run /var/lib/clixon

[Install]
WantedBy=multi-user.target
