# build-pkg/clixon-vpp_1.1.0_amd64/DEBIAN/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Create clixon group if not exists
        if ! getent group clixon > /dev/null 2>&1; then
            groupadd -r clixon
        fi
        
        # Create clixon user if not exists  
        if ! getent passwd clixon > /dev/null 2>&1; then
            useradd -r -g clixon -d /var/lib/clixon -s /usr/sbin/nologin -c "Clixon daemon" clixon 2>/dev/null || true
        fi
        
        # Add www-data to clixon group
        usermod -a -G clixon www-data 2>/dev/null || true
        
        # Set permissions
        chown -R root:clixon /var/lib/clixon 2>/dev/null || true
        chmod 775 /var/lib/clixon/vpp 2>/dev/null || true
        
        # Update library cache
        ldconfig
        
        # Reload systemd
        systemctl daemon-reload || true
        
        echo ""
        echo "=============================================="
        echo "  clixon-vpp installed successfully!"
        echo "=============================================="
        echo ""
        echo "To start services:"
        echo "  sudo systemctl start clixon-vpp-backend"
        echo "  sudo systemctl start clixon-vpp-restconf"
        echo ""
        echo "To enable on boot:"
        echo "  sudo systemctl enable clixon-vpp-backend clixon-vpp-restconf"
        echo ""
        echo "CLI (use either):"
        echo "  cli"
        echo "  clixon_cli -f /etc/clixon/clixon-vpp.xml"
        echo ""
        echo "RESTCONF:"
        echo "  http://localhost:8080/restconf/data"
        echo ""
        ;;
esac
exit 0
